# movie.socket
[Unit]
Description=movie socket
[Socket]
ListenStream=/run/movie.sock
[Install]
WantedBy=sockets.target
